# daemon/systemd/mailinabox-helper.service
# Privileged helper for the Mail-in-a-Box management daemon.
# Executes a fixed menu of privileged operations on behalf of the
# unprivileged manager over /run/mailinabox/helper.sock.
#
# Installed to /etc/systemd/system/ by setup; the binary is installed to
# /usr/local/lib/mailinabox/ (never run from the repo checkout).

[Unit]
Description=Mail-in-a-Box privileged helper
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/lib/mailinabox/helperd -socket /run/mailinabox/helper.sock -socket-group mailinabox
RuntimeDirectory=mailinabox
RuntimeDirectoryMode=0750
Restart=on-failure
RestartSec=2

# Hardening. ProtectSystem is deliberately NOT set: host.apt_upgrade
# writes across /usr and /etc via dpkg. Revisit if apt moves to a
# transient systemd-run unit (open question in helper-intent-menu.md).
NoNewPrivileges=yes
PrivateTmp=yes
ProtectHome=read-only
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes

[Install]
WantedBy=multi-user.target
